# f2_control/files/general/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh
# no-op: let s6 restart nginx if it exits
